# cookbooks/letsencrypt/templates/crontab/letsencrypt.erb
#!/bin/sh

. /etc/environment && \
    certbot renew -n \
        1>/var/log/letsencrypt.log \
        2>/var/log/letsencrypt.err && \
    nginx -s reload
