# cookbooks/znc/templates/systemd/znc.service.erb
[Unit]
Description = ZNC service

[Service]
User = <%= @user %>
Group = <%= @group %>
WorkingDirectory = <%= @working_directory %>
ExecStart = /usr/bin/znc -f
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
