# cookbooks/rakudobrew/templates/etc/profile.d/rakudobrew.sh.erb
export PATH=<%= @install_dir %>/bin:$PATH
eval "$(<%= @install_dir %>/bin/rakudobrew init -)"
