# cookbooks/anyenv/templates/etc/profile.d/anyenvrc.sh.erb
export ANYENV_ROOT=<%= @install_dir %>
export PATH="<%= @install_dir %>/bin:$PATH"
eval "$(anyenv init -)"
