# cookbooks/helloworld/templates/helloworld.service.erb
[Unit]
Description = Service for Hello World

[Service]
User = <%= @user %>
Group = <%= @group %>
WorkingDirectory = <%= @workdir %>
ExecStart = <%= @workdir %>/bin/helloworld
Restart = always
RestartSec = 30s
Type = simple

[Install]
WantedBy = multi-user.target
